# local_mac/launchd/com.amitgrupper.klaus.things-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amitgrupper.klaus.things-poller</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/amitgrupper/Desktop/Klaus/.venv/bin/python</string>
        <string>-m</string>
        <string>local_mac.things_poller</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amitgrupper/Desktop/Klaus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- WHY load_dotenv: things_poller.py calls load_dotenv(override=True) from WorkingDirectory.
         All GCP_PROJECT_ID, FIRESTORE_* vars are loaded from .env automatically.
         Only PATH is set here to help subprocesses (osascript) find system binaries. -->

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/amitgrupper/Library/Logs/klaus-things-poller.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/amitgrupper/Library/Logs/klaus-things-poller.log</string>
</dict>
</plist>
